# com.kian.transcribe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        This launchd plist sets up a service with the label 'com.kian.transcribe'.
        It does the following:

        - Watches these two iCloud directories:
            /Users/kian/Library/Mobile Documents/com~apple~CloudDocs/Music/notes
            /Users/kian/Library/Mobile Documents/com~apple~CloudDocs/Music/course
        - When a change is detected in either directory, it runs:
            /Users/kian/Developer/transcribe/run_transcribe.sh
        - The script runs with the working directory:
            /Users/kian/Developer/transcribe
    -->
    <key>Label</key>
    <string>com.kian.transcribe</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/kian/Developer/transcribe/run_transcribe.sh</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>/Users/kian/Library/Mobile Documents/com~apple~CloudDocs/Music/notes</string>
        <string>/Users/kian/Library/Mobile Documents/com~apple~CloudDocs/Music/course</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/kian/Developer/transcribe</string>
</dict>
</plist>
